# systemd/nanobot-rag-ingest.service.template
[Unit]
Description=nanobot RAG bridge periodic ingestion
After=nanobot-rag.service

[Service]
Type=oneshot
ExecStart=${LOCAL_BIN_DIR}/nanobot-rag-ingest
